# configfiles/opwifi
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
	/sbin/wifi up
}

stop() {
	/sbin/wifi down
}

restart() {
	stop
	start
}
