# configfiles/opwifi
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
	/sbin/wifi up
	return 0
}

stop() {
	/sbin/wifi down
	return 0
}

restart() {
    if stop; then
        start
    fi
}
